# roles/xforsa.miners/templates/xorg.service.j2
[Unit]
Description=overclocking gpus
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/startx

[Install]
WantedBy=default.target
